# server/src/main/deb/DEBIAN/prerm
#!/bin/bash

systemctl stop 'soffid-iamsync.service' || true

systemctl disable 'soffid-iamsync.service' || true
